# overlay/etc/init.d/10_usb_gadget
#!/sbin/openrc-run

start() {
    /usr/bin/usb_gadget start || true
    /usr/bin/unudhcpd -i usb0 -s 172.16.42.1 -c 172.16.42.2 &
}

stop() {
    # dont
    true
}